# services/keepalived/run
#!/usr/bin/with-contenv bash
exec /usr/sbin/keepalived -l -G -f /etc/pihole/keepalived.conf --dont-fork
