# services/stubby/run
#!/usr/bin/with-contenv bash
exec s6-setuidgid nobody /usr/bin/stubby -C /etc/pihole/stubby.yml
