# services/stubby/run
#!/usr/bin/with-contenv bash
exec s6-setuidgid stubby /usr/bin/stubby -C /etc/pihole/stubby.yml
